# config/systemd/user/check-battery.timer
[Unit]
Description=Timer to check battery status

[Timer]
OnActiveSec=3min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target