# config/systemd/user/check-battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/usr/bin/check-battery