# config/systemd/user/kite-updater.service
[Unit]
Description=Kite Updater

[Service]
Type=oneshot
ExecStart=/home/antcc/.local/share/kite/current/kite-update "update"

[Install]
WantedBy=multi-user.target