# config/systemd/user/list-pkgs.service
[Unit]
Description=list-pkgs

[Service]
ExecStart=/usr/bin/atom-list-pkgs
ExecStart=/usr/bin/pacman-list-pkgs