# config/systemd/user/update-dots.timer
[Unit]
Description=Timer to update dotfiles

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WandedBy=timers.target
