# config/systemd/user/update-dots.service
[Unit]
Description=Update dotfiles repo

[Service]
Type=oneshot
ExecStart=/usr/bin/update-dots
