# etc/systemd/system/gdm.service.d/override.conf
[Service]
ExecStartPre=~/.dotfiles/bin/monitor-hotplug.sh
